# minecraft.service
[Unit]
Description=Minecraft Server (Adrenaline Modpack)
After=network.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=forking
User=minecraft
Group=minecraft
WorkingDirectory=/opt/minecraft-server

# Variável de ambiente do terminal e screen
Environment="TERM=xterm"
Environment="SCREENDIR=/opt/minecraft-server/.screen"

# Garante que as pastas necessárias existam e usa o lugar certo pro screen
ExecStartPre=+/usr/bin/mkdir -p /opt/minecraft-server/.screen
ExecStartPre=+/usr/bin/chown minecraft:minecraft /opt/minecraft-server/.screen

# Comando de inicialização usando screen
ExecStart=/usr/bin/screen -dmS minecraft /opt/minecraft-server/start-server.sh

# Comando de parada (usa o script de gerenciamento)
ExecStop=/opt/minecraft-server/mc-manager.sh stop

# Timeout para parada
TimeoutStopSec=60

# Reiniciar em caso de falha
Restart=on-failure
RestartSec=30

# Limites de recursos
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=3200M

# Prioridade de CPU e I/O (menor número = maior prioridade)
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=0

# OOM killer - não matar o servidor facilmente
OOMScoreAdjust=-800

# Ambiente
Environment="JAVA_HOME=/usr/lib/jvm/java-21-openjdk"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Diretórios protegidos
ProtectSystem=no
ProtectHome=read-only
ReadWritePaths=/opt/minecraft-server

[Install]
WantedBy=multi-user.target
