# minecraft/minecraft.service
[Unit]
Description=Minecraft Server
After=network.target
StartLimitIntervalSec=60s
StartLimitBurst=3

[Service]
Type=simple
User=__SERVER_USER__
Group=__SERVER_USER__
WorkingDirectory=__SERVER_DIR__
ExecStart=__SERVER_DIR__/start-server.sh
# Prefer SIGINT to give the JVM a chance to shutdown cleanly.
KillSignal=SIGINT
TimeoutStopSec=120
Restart=on-failure
RestartSec=15
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=__MEMORY_MAX_MB__M
Nice=-5
IOSchedulingClass=best-effort
IOSchedulingPriority=0
OOMScoreAdjust=-800
ProtectSystem=no
ProtectHome=read-only
ReadWritePaths=__SERVER_DIR__

[Install]
WantedBy=multi-user.target
